# etc/rgbmatrix.service
[Unit]
Description=RGB Matrix
Requires=xvfb.service ft-server.service
After=xvfb.service ft-server.service

[Service]
WorkingDirectory=/opt/rgbmatrix
ExecStart=/bin/bash ./start.sh
Environment=PYTHONUNBUFFERED=1
Restart=always
Type=simple
StandardOutput=null
StandardError=null

[Install]
WantedBy=default.target
